# etc/init.d/acpid-openrc
#!/sbin/openrc-run

description="Starts the Advanced Configuration and Power Interface daemon"
command="/usr/sbin/acpid"
command_args="-f ${ACPI_OPTS}"
supervisor=supervise-daemon

depend() {
	need udev
	after modules
	provide acpid
}

start_pre() {
	if [ -n "${ACPI_MODULES}" ]; then
		ebegin "Loading required ACPI modules"
		modprobe -a "${ACPI_MODULES}"
		eend $?
	fi
}
