# etc/init.d/acpi-support-openrc
#!/sbin/openrc-run

description="Applies ACPI power management settings"
ACPI_FUNCTIONS="/usr/share/acpi-support/power-funcs"

depend() {
	after local
	need acpid
}

start_pre() {
	if [ ! -f "${ACPI_FUNCTIONS}" ]; then
		eerror "ACPI support functions not found at ${ACPI_FUNCTIONS}"
		return 1
	fi
	# shellcheck source=/dev/null
	. "${ACPI_FUNCTIONS}"
}

start() {
	ebegin "Applying ACPI power management settings"

	if on_ac_power; then
		einfo "On AC power, no changes needed."
	else
		einfo "On battery power, applying power saving settings."
		/etc/acpi/power.sh true
	fi

	# Source all startup scripts.
	if [ -d /etc/acpi/start.d ]; then
		for SCRIPT in /etc/acpi/start.d/*.sh; do
			[ -f "${SCRIPT}" ] && . "${SCRIPT}"
		done
	fi

	eend 0
}

stop() {
	if [ ! -f "${ACPI_FUNCTIONS}" ]; then
		return 0
	fi
	# shellcheck source=/dev/null
	. "${ACPI_FUNCTIONS}"

	ebegin "Restoring ACPI power settings"
	if ! on_ac_power; then
		/etc/acpi/power.sh false
	fi
	eend 0
}
